# delay.service
[Unit]
Description=Delay

[Service]
ExecStart=/root/delay

[Install]
WantedBy=multi-user.target
